# dockerfiles/supervisord.conf
[supervisord]
nodaemon=true

[program:nodejs]
autostart=true
user=root
command=/usr/bin/nodejs /usr/bin/reveal-md /tmp/index.md --theme chef
